# supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
command=/usr/bin/redis-server

[program:cosnok]
command=/usr/local/bundle/bin/bundle exec ruby bot.rb
